# systemd/archive-data.service
[Unit]
Description=Archive Polymarket Old Data (30+ days)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/trader-data-collect
ExecStart=/home/ubuntu/trader-data-collect/venv/bin/python3 /home/ubuntu/trader-data-collect/archive_old_data.py --days 30 --hot-dir /home/ubuntu/trader-data-collect/real_hot --archive-dir /home/ubuntu/trader-data-collect/real_archive
StandardOutput=append:/home/ubuntu/trader-data-collect/logs/archive.log
StandardError=append:/home/ubuntu/trader-data-collect/logs/archive.log

# 资源限制
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=default.target

