# systemd/cex-score-1h.service
[Unit]
Description=CEX Score Daemon (1h)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/trader-data-collect
ExecStart=/home/ubuntu/trader-data-collect/venv/bin/python3 /home/ubuntu/trader-data-collect/cex_score_daemon.py --symbol btc --window 1h --port 9002 --hot-dir /home/ubuntu/trader-data-collect/real_hot --lookback-s 21600 --decay-T 60 --decay-lambda-base 0.02 --decay-sigma 6 --decay-multiplier 0.5 --decay-min-mu 8 --decay-max-mu 60 --decay-N-windows 10

# 自动重启
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cex-score-1h

# 资源限制（使用新语法）
MemoryMax=512M
CPUQuota=50%

# 环境变量
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
