# systemd/cex-recorder.service
[Unit]
Description=CEX Multi-Asset Order Book Recorder
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trader-data-collect/scripts
ExecStart=/home/ubuntu/trader-data-collect/venv/bin/python3 /home/ubuntu/trader-data-collect/scripts/cex_multi_asset_recorder.py

# 自动重启
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cex-recorder

# 资源限制
MemoryLimit=512M
CPUQuota=50%

# 环境变量
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

